# etc/crontabs/root
*/1 * * * * /root/aquarium.sh
1 1 * * * /root/createSunInfo.sh
